# patch/files/etc/uci-defaults/90-container-ports
#!/bin/sh
# ============================================
# 自动放行容器端口（Podman/Docker）
# ============================================

# 需要放行的容器端口列表
CONTAINER_PORTS="5908 3000 8080 9090 8443 5432 6379 9200 2379 2380 8888 6080 1566"

for port in $CONTAINER_PORTS; do
  # 检查是否已存在
  if ! uci get firewall."container_ports_$port".name >/dev/null 2>&1; then
    uci set firewall."container_ports_$port"="rule"
    uci set firewall."container_ports_$port".name="Allow Container Port $port"
    uci set firewall."container_ports_$port".proto="tcp"
    uci set firewall."container_ports_$port".dest_port="$port"
    uci set firewall."container_ports_$port".target="ACCEPT"
    uci set firewall."container_ports_$port".family="ipv4"
  fi
done

# 同时放行 IPv6
for port in $CONTAINER_PORTS; do
  if ! uci get firewall."container_ports_ipv6_$port".name >/dev/null 2>&1; then
    uci set firewall."container_ports_ipv6_$port"="rule"
    uci set firewall."container_ports_ipv6_$port".name="Allow Container Port $port (IPv6)"
    uci set firewall."container_ports_ipv6_$port".proto="tcp"
    uci set firewall."container_ports_ipv6_$port".dest_port="$port"
    uci set firewall."container_ports_ipv6_$port".target="ACCEPT"
    uci set firewall."container_ports_ipv6_$port".family="ipv6"
  fi
done

uci commit firewall
