# patch/files/etc/init.d/podman-daemon
#!/bin/sh /etc/rc.common

START=99
USE_PROCD=1
NAME=podman-daemon
PROG=/usr/bin/podman

start_service() {
 mkdir -p /run/podman
 procd_open_instance
 procd_set_param command "$PROG" system service -t 0
 procd_set_param respawn
 procd_set_param file /run/podman/podman.sock
 procd_set_param stdout 1
 procd_set_param stderr 1
 procd_close_instance

 # 自动扫描并放行容器端口
 /etc/init.d/podman-daemon collect_ports
}

stop_service() {
 killall podman 2>/dev/null
}

collect_ports() {
 # 扫描所有运行中容器的端口映射
 local ports=""
 for port_line in $(podman ps --format '{{.Names}} {{.Ports}}' 2>/dev/null | grep -oP '\d+/tcp' | grep -oP '^\d+'); do
   ports="$ports $port_line"
 done

 # 去重
 ports=$(echo "$ports" | tr ' ' '\n' | sort -un | tr '\n' ' ')

 for port in $ports; do
   # 检查是否已放行
   if ! iptables -C INPUT -p tcp --dport "$port" -j ACCEPT 2>/dev/null; then
     iptables -I INPUT -p tcp --dport "$port" -j ACCEPT
     ip6tables -I INPUT -p tcp --dport "$port" -j ACCEPT 2>/dev/null
     logger -t podman-daemon "Opened port $port for container"
   fi
 done
}
